# systemd/gcal-bridge.service
# Example systemd unit (user or system)
# - If you run as a *user* service: `systemctl --user enable --now gcal-bridge`
# - If you run as a *system* service: put in /etc/systemd/system and use sudo.

[Unit]
Description=GCal/Gmail bridge (UDS JSONL) for lobs-mcp
After=network.target

[Service]
Type=simple

# Run the bridge under a dedicated user; lobs-mcp (run by rafe) connects via a shared group.
# You must:
#   sudo groupadd -f lobs-mcp
#   sudo usermod -aG lobs-mcp rafe
#   sudo useradd -m -s /bin/bash mpcuser   (if not already)
#   sudo usermod -aG lobs-mcp mpcuser
User=mpcuser
Group=lobs-mcp
UMask=0007

# Create /run/lobs-mcp with safe permissions for the socket.
RuntimeDirectory=lobs-mcp
RuntimeDirectoryMode=0770

# Bridge code location:
# Since this service runs as `mpcuser`, the simplest approach is: keep a copy of the repo in
#   /home/mpcuser/lobs-mcp
WorkingDirectory=%h/lobs-mcp

Environment=LOBS_BRIDGE_SOCKET=/run/lobs-mcp/bridge.sock
Environment=LOBS_BRIDGE_TIMEOUT_MS=15000
Environment=LOBS_BRIDGE_LOG_LEVEL=info
Environment=LOBS_BRIDGE_SOCKET_MODE=0660

ExecStart=%h/lobs-mcp/bin/start-bridge
Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target
