# systemd/gcal-bridge.user.service
# User service (recommended): keeps bridge alive and restartable.
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp systemd/gcal-bridge.user.service ~/.config/systemd/user/gcal-bridge.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now gcal-bridge
# Logs:
#   journalctl --user -u gcal-bridge -f

[Unit]
Description=GCal/Gmail bridge (UDS JSONL) for lobs-mcp
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/lobs-mcp
Environment=LOBS_BRIDGE_SOCKET=/run/lobs-mcp/bridge.sock
Environment=LOBS_BRIDGE_TIMEOUT_MS=15000
Environment=LOBS_BRIDGE_LOG_LEVEL=info
Environment=LOBS_BRIDGE_SOCKET_MODE=0660
ExecStart=%h/lobs-mcp/bin/start-bridge
# (start-bridge is build-smart; it won't rebuild unless dist is missing or LOBS_FORCE_BUILD=1)
Restart=always
RestartSec=1

[Install]
WantedBy=default.target
